# borg-exporter.service
[Unit]
Description=Prometheus Borg Exporter
After=network-online.target

[Service]
ExecStart=/usr/local/bin/borg-exporter.sh
Type=oneshot