# borg-exporter.timer
[Unit]
Description=Prometheus Borg Exporter Timer

[Timer]
OnCalendar=hourly

[Install]
WantedBy=timers.target